# agents/ac/old/ac.service
[Unit]
Description=ac
After=network-online.target network.target mosquitto.service influxdb.service

[Service]
Type=simple
EnvironmentFile=/etc/sd.env
ExecStart=/opt/sd/bin/ac -c /opt/sd/etc/ac.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
